# pi/services/home-center-kiosk.service
[Unit]
Description=Home Center Chromium Kiosk
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/__USER__/.Xauthority
# For Wayland (Pi OS Bookworm default)
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=__REPO_DIR__/pi/kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
