# pi/services/mic-streamer.service
[Unit]
Description=Home Center XVF3800 mic streamer
After=sound.target
Wants=sound.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
ExecStart=__REPO_DIR__/pi/.venv/bin/python __REPO_DIR__/pi/mic_streamer.py
Restart=always
RestartSec=3
SupplementaryGroups=audio
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mic-streamer

[Install]
WantedBy=multi-user.target
