# pi/services/morning-dashboard.service
[Unit]
Description=Home Center weekday morning dashboard HDMI activation
Wants=dashboard-local.service
After=dashboard-local.service

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__REPO_DIR__
ExecStart=__REPO_DIR__/pi/morning_dashboard.sh

