# pi/services/wake-word.service
[Unit]
Description=Home Center command server (dashboard state, CEC, timers, chime)
After=sound.target
Wants=sound.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
# --no-wake-detection: the wake-word DNN + Whisper now run on the Mac mini
# voice-service and drive this server via HTTP (/api/tv/on, /api/chime, …).
ExecStart=__REPO_DIR__/pi/.venv/bin/python __REPO_DIR__/pi/wake_word_service.py --no-wake-detection
Restart=always
RestartSec=3

# Audio device access
SupplementaryGroups=audio video

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wake-word

[Install]
WantedBy=multi-user.target
