# docs/examples/com.homecenter.design-claw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template launchd job for the Home Center Design Claw daily run.

  NOT INSTALLED. This file is documentation only.

  To use:
    1. Copy this file.
    2. Replace placeholders:
         __REPO_DIR__              e.g. /Users/peter/home-center
         __PYTHON_BIN__            e.g. /opt/homebrew/bin/python3 or a venv's python3
         __OPENAI_API_KEY__        your OpenAI API key
         __TELEGRAM_BOT_TOKEN__    Design Claw bot token (NOT the OpenClaw family bot)
         __TELEGRAM_CHAT_ID__      target chat (e.g. your personal DM, not the family group)
    3. Move to ~/Library/LaunchAgents/com.homecenter.design-claw.plist
    4. launchctl load ~/Library/LaunchAgents/com.homecenter.design-claw.plist

  Schedule: 08:15 local, daily. Edit <StartCalendarInterval> if desired.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homecenter.design-claw</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>__REPO_DIR__/scripts/run_daily_design_and_send.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>__TELEGRAM_BOT_TOKEN__</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>__TELEGRAM_CHAT_ID__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/design_outputs/.launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/design_outputs/.launchd.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
