# deploy/mac-mini/com.homecenter.voice-homecenter2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homecenter.voice.homecenter2</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO__/voice-service/.venv/bin/python3</string>
        <string>__REPO__/voice-service/voice_service.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO__/voice-service</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>MIC_STREAM_HOST</key>
        <string>homecenter2.local</string>
        <key>MIC_STREAM_PORT</key>
        <string>8766</string>
        <key>PI_COMMAND_URL</key>
        <string>http://homecenter2.local:8765</string>
        <key>WORKER_URL</key>
        <string>https://home-center-api.phhowell.workers.dev</string>
        <key>WORKER_TOKEN</key>
        <string>__WORKER_TOKEN__</string>
        <key>OPENCLAW_BRIDGE_URL</key>
        <string>http://127.0.0.1:3100</string>
        <key>OPENCLAW_CHAT_ID</key>
        <string>__OPENCLAW_CHAT_ID__</string>
        <key>DESIGN_FEEDBACK_SCRIPT</key>
        <string>__REPO__/scripts/record_design_feedback.py</string>
        <key>VOSK_MODEL_DIR</key>
        <string>__REPO__/voice-service/models/vosk-model-small-en-us-0.15</string>
        <key>WHISPER_MODEL</key>
        <string>base.en</string>
        <key>OPENAI_STT_MODE</key>
        <string>fallback</string>
        <key>OPENAI_STT_MODEL</key>
        <string>gpt-4o-transcribe</string>
        <key>VOICE_RELIABILITY_LOG</key>
        <string>__REPO__/voice-service/logs/voice-reliability-homecenter2.jsonl</string>
        <key>POST_ACTION_MUTE_SECONDS</key>
        <string>0.15</string>
        <key>WAKE_COOLDOWN</key>
        <string>0.75</string>
        <key>PRE_WAKE_SECONDS</key>
        <string>2.4</string>
        <key>COMMAND_MIN_SECONDS</key>
        <string>0.16</string>
        <key>COMMAND_MAX_SECONDS</key>
        <string>0.9</string>
        <key>POST_WAKE_SILENCE_SECONDS</key>
        <string>0.08</string>
        <key>FOLLOWUP_COMMAND_WAIT_SECONDS</key>
        <string>2.0</string>
        <key>FOLLOWUP_COMMAND_MAX_SECONDS</key>
        <string>3.0</string>
        <key>WAKE_ENGINE</key>
        <string>vosk</string>
        <key>WAKE_CONFIRM_COMMAND</key>
        <string>0</string>
        <key>CONFIRM_PRE_WAKE_SECONDS</key>
        <string>5.0</string>
        <key>CONFIRM_POST_WAKE_SECONDS</key>
        <string>2.0</string>
        <key>OPENWAKEWORD_MODEL</key>
        <string>__REPO__/pi/models/hey_homer.onnx</string>
        <key>OPENWAKEWORD_THRESHOLD</key>
        <string>0.92</string>
        <key>OPENWAKEWORD_MIN_CONSECUTIVE</key>
        <string>3</string>
        <key>OPENWAKEWORD_AUDIO_HEARTBEAT_SECONDS</key>
        <string>2.0</string>
        <key>OPENWAKEWORD_AUDIO_LOG_MIN_RMS</key>
        <string>180</string>
        <key>OPENWAKEWORD_AUDIO_LOG_INTERVAL_SECONDS</key>
        <string>1.0</string>
        <key>SPEECH_CANDIDATE_MIN_PEAK_RMS</key>
        <string>450</string>
        <key>SPEECH_CANDIDATE_END_SILENCE_SECONDS</key>
        <string>0.65</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO__/voice-service/logs/voice-homecenter2-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/voice-service/logs/voice-homecenter2-stderr.log</string>
</dict>
</plist>
