# deploy/mac-mini/com.homecenter.takeout-suggestions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homecenter.takeout-suggestions</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__REPO__/scripts/update-takeout-suggestions.py</string>
    <string>--days</string>
    <string>90</string>
    <string>--max</string>
    <string>50</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TAKEOUT_GMAIL_ACCOUNTS</key>
    <string>__TAKEOUT_GMAIL_ACCOUNTS__</string>
    <key>WORKER_URL</key>
    <string>https://home-center-api.phhowell.workers.dev</string>
    <key>WORKER_TOKEN</key>
    <string>__WORKER_TOKEN__</string>
  </dict>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__REPO__/logs/takeout-suggestions.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/logs/takeout-suggestions.err.log</string>
</dict>
</plist>
