# deploy/mac-mini/com.homecenter.voice-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homecenter.voice-health</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO__/voice-service/.venv/bin/python3</string>
        <string>__REPO__/voice-service/health_check.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO__/voice-service</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PI_COMMAND_URL</key>
        <string>http://homecenter.local:8765</string>
        <key>VOICE_HEALTH_MAX_LOG_AGE_SECONDS</key>
        <string>3600</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>15</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>45</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO__/voice-service/logs/voice-health.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/voice-service/logs/voice-health.err.log</string>
</dict>
</plist>
