# deploy/mac-mini/com.homecenter.design-claw-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template for the Home Center Design Claw Telegram listener launchd
  agent — pairs with com.homecenter.design-claw.plist (the daily job).

  Rendered + loaded by deploy/mac-mini/setup-design-claw.sh. Do NOT
  commit a copy with real values.

    __REPO_DIR__              absolute path to the home-center clone
    __PYTHON_BIN__            full path to python3
    __OPENAI_API_KEY__        OpenAI API key (same one the daily uses)
    __TELEGRAM_BOT_TOKEN__    Design Claw bot token (same as daily)
    __TELEGRAM_CHAT_ID__      Target DM chat id
    __DESIGN_CLAW_MODEL__     Model id (default gpt-5.4-mini)

  Polls every 5 min at the :00/:05/:10/.../:55 mark. Edit the
  <StartCalendarInterval> array to change cadence.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homecenter.design-claw-listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>__REPO_DIR__/scripts/design_listener.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>__TELEGRAM_BOT_TOKEN__</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>__TELEGRAM_CHAT_ID__</string>
        <key>DESIGN_CLAW_MODEL</key>
        <string>__DESIGN_CLAW_MODEL__</string>
        <key>DESIGN_CLAW_IMAGE_MODEL</key>
        <string>__DESIGN_CLAW_IMAGE_MODEL__</string>
    </dict>

    <!-- Calendar-based schedule fires more reliably than StartInterval
         in GUI-session LaunchAgents (which can get paused by AppNap /
         QoS throttling). Every 5 min at :00, :05, :10, ..., :55. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>5</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>15</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>25</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>35</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>45</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
        <dict><key>Minute</key><integer>55</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/design_outputs/.launchd-listener.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/design_outputs/.launchd-listener.stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
