# .openclaw/mac-mini/plists/com.openclaw.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>__REPO_DIR__/openclaw/index.js</string>
        <string>--port</string>
        <string>3100</string>
        <string>--llm-url</string>
        <string>https://home-center-api.phhowell.workers.dev/api/ask</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__/openclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <!-- Set TELEGRAM_BOT_TOKEN here (from @BotFather) before loading this plist. -->
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>__TELEGRAM_BOT_TOKEN__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/.openclaw/logs/openclaw-bridge-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/.openclaw/logs/openclaw-bridge-stderr.log</string>
</dict>
</plist>
