# openclaw/openclaw.service
[Unit]
Description=OpenClaw Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/home-center/openclaw
ExecStart=/usr/bin/node index.js --port 3100 --llm-url https://home-center-api.phhowell.workers.dev/api/ask
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=-/home/pi/home-center/openclaw/.env

[Install]
WantedBy=multi-user.target
